# sendy.crontab
# Scheduling Campaigns
*/5 * * * * curl -fsSL http://localhost:8080/scheduled.php

# Autoresponders
*/1 * * * * curl -fsSL http://localhost:8080/autoresponders.php

# CSV Importer
*/1 * * * * curl -fsSL http://localhost:8080/import-csv.php