# s6-rc.d/supercronic/run
#!/bin/execlineb -P
# S6 service to run supercronic for scheduled tasks
# This service will be automatically restarted by S6 if it exits

if -t { test -f /etc/sendy.crontab }
if -t { /usr/local/bin/supercronic -test /etc/sendy.crontab }
fdmove -c 2 1
exec /usr/local/bin/supercronic -debug /etc/sendy.crontab
